# meta-evb/meta-evb-nuvoton/meta-buv-runbmc/recipes-phosphor/gpio/phosphor-gpio-monitor/psu_hwmon_stop.service
[Unit]
Description=hwmon stop

[Service]
Type=oneshot
ExecStart=/bin/sh -c "psu_hotplug_action.sh 2"
SyslogIdentifier=hwmon_stop
