# meta-evb/meta-evb-nuvoton/meta-buv-runbmc/recipes-phosphor/gpio/phosphor-gpio-monitor/psu_hwmon_start.service
[Unit]
Description=hwmon start

[Service]
Type=oneshot
ExecStart=/bin/sh -c "psu_hotplug_action.sh 3"
SyslogIdentifier=hwmon_start
