# meta-nuvoton-obmc/meta-evb-npcm845/recipes-evb-npcm845/usb-tty/files/usb_tty.service
[Unit]
Description=enable usb tty

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/sbin/usb_tty.sh'

[Install]
WantedBy=multi-user.target
