# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-phosphor/gpio/phosphor-gpio-monitor/dp_hotplug.service
[Unit]
Description=dp hotplug detected

[Service]
Type=oneshot
ExecStart=/usr/bin/dp_hotplug.sh
SyslogIdentifier=dp_hotplug

