# meta-nuvoton-obmc/meta-common/recipes-nuvoton/nist-linux-selftest/files/nist-linux-selftest.service
[Unit]
Description=nist-linux-selftest Service

[Service]
Type=oneshot
ExecStart=/usr/sbin/nist-linux-selftest

[Install]
WantedBy=multi-user.target
